# archy-executor.service
[Unit]
Description=Archy Executor Daemon - AI System Assistant
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/home/chef/Archy/target/release/archy-executor

# Auto-restart configuration
# Restart on any failure or process exit
Restart=always
RestartSec=2

# Prevent restart storms - allow 5 restarts within 60 seconds
StartLimitInterval=60
StartLimitBurst=5
StartLimitAction=none

# Resource limits to prevent memory exhaustion
MemoryMax=512M

# Process management
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=5

# Logging to systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=archy-executor


# Working directory
WorkingDirectory=/home/chef/Archy

# Environment variables
Environment="RUST_LOG=info"


[Install]
# Start at boot (user-level)
WantedBy=default.target


